# tools/flipper-badge-bridge/flipper-badge-bridge.service
[Unit]
Description=DC32 Flipper QFlipper RPC badge bridge
After=local-fs.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/flipper-badge-bridge
ExecStart=/opt/flipper-badge-bridge/flipper_badge_bridge.py
Restart=always
RestartSec=2
PrivateNetwork=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
